# systemd/upgrade-system.service
[Unit]
Description=Upgrade all system

[Service]
Type=oneshot
ExecStart=/usr/bin/pacman -Syu -q --noconfirm
