# conf/daily_update/crontab
* 11 * * * /usr/bin/pacman -Syu -q --noconfirm
