# conf/daily_update/crontab
0 11 * * * root /usr/bin/pacman -Syu -q --noconfirm
